# dist/debian/scrimplb.service
[Unit]
Description=Scrimp load balancer
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/scrimplb -config-file /etc/scrimplb/scrimp-lb.json
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=scrimplb
User=scrimplb
Group=scrimplb
Restart=on-failure
RestartSec=10s
StartLimitInterval=1min

[Install]
WantedBy=multi-user.target
